# crontab
MAILTO=test@test.test

0   *   *   *   *    cd /dir && ./script_hw7.sh
